# systemd/windfarm-monitoring.timer
[Unit]
Description=Restart windfarm-monitoring

[Timer]
OnUnitInactiveSec=10min

[Install]
WantedBy=timers.target